# contrib/sans-car-bot.service
[Unit]
Description=Sans car Discord bot
Requires=lavalink.service postgresql.service

[Service]
# Replace with your system config
WorkingDirectory=/opt/sans-car/
ExecStart=/opt/sans-car/deployable/Bot
Restart=always
# Restart service after 10 seconds if the dotnet service crashes:
RestartSec=10
KillSignal=SIGINT
SyslogIdentifier=sans-bot
User=sanscar
Environment=DOTNET_PRINT_TELEMETRY_MESSAGE=false
Environment=DOTNET_CLI_TELEMETRY_OPTOUT=1

[Install]
WantedBy=multi-user.target
